# etc/sv/clipboard_backup/finish
#!/bin/sh

. "$(dirname "$0")/utils.sh"

[ -d "$CLIPBOARD_BACKUP_FOLDER" ] && {
	for file in "$CLIPBOARD_BACKUP_FOLDER"/*; do
		shred -uz "$file"
	done
	rm -rf "$CLIPBOARD_BACKUP_FOLDER"
}
